# polier_pro/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
bashio::log.error "nginx exited with code $1"
