# polier_pro/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting nginx on port 3000..."
nginx -t && bashio::log.info "nginx config OK" || bashio::log.error "nginx config FAILED"
exec nginx -g "daemon off;" 2>&1
